# dot_config/systemd/user/autolock.timer
[Unit]
Description=Lock if user is idle

[Timer]
OnCalendar=minutely

[Install]
WantedBy=timers.target
